# docker/supervisord-apache2.conf
[program:apache2]
command=/utils/start-apache2.sh
numprocs=1
autostart=true
autorestart=true